# docker/kadserver.sv.ini
[program:kadserver]
command=/usr/local/bin/node /home/flunky/kadserver/index.js
priority=1
user=flunky
stdout_logfile=/var/log/kadserver.log
stderr_logfile=/var/log/kadserver.err
