# docker/nodejs.sv.ini
[program:nodejs]
command=/usr/local/bin/node /home/nodejs/src/kadserver.js
priority=1
user=nodejs
stdout_logfile=/var/log/nodejs.log
stderr_logfile=/var/log/nodejs.err
